# systemd/system/aptupdatechecker_autoremove.service
[Unit]
Description=Unnecessary package autoremover

[Service]
Type=oneshot
CPUQuota=25%
ExecStart=apt-get -y autoremove
