# systemd/system/aptupdatechecker_autoclean.service
[Unit]
Description=Package cache autocleaner

[Service]
Type=oneshot
CPUQuota=25%
ExecStart=apt-get -y autoclean

