# services/my_apt_autoclean.service
[Unit]
Description=Package cache autocleaner

[Service]
Type=oneshot
CPUQuota=50%
ExecStart=apt-get -y autoclean
