# services/my_apt_update_checker.service
[Unit]
Description=Checks for package updates
After=network-online.target

[Service]
Type=oneshot
CPUQuota=50%
ExecStart=