# services/my_firmware_update_checker.service
[Unit]
Description=Checks for firmware updates
After=network-online.target

[Service]
Type=oneshot
CPUQuota=50%
ExecStart=