# systemd/aptupdatechecker_apt_update_checker.service
[Unit]
Description=Checks for package updates

[Service]
Type=oneshot
CPUQuota=50%
ExecStart=/usr/bin/aptupdatechecker apt
