# systemd/aptupdatechecker_firmware_update_checker.service
[Unit]
Description=Checks for firmware updates

[Service]
Type=oneshot
CPUQuota=50%
ExecStart=/usr/bin/aptupdatechecker firmware
