# debian/aptupdatechecker_autoremove.service
[Unit]
Description=Unnecessary package autoremover

[Service]
Type=oneshot
CPUQuota=50%
ExecStart=apt-get -y autoremove
