# debian/postinst
#!/bin/sh
set -e

#DEBHELPER#

if [ "$1" = "configure" ]; then
    # Reload systemd daemon to pick up new unit files
    deb-systemd-helper daemon-reload || true

    # Enable the timers
    deb-systemd-helper enable aptupdatechecker.service || true
    deb-systemd-helper enable aptupdatechecker_apt_update_checker.timer || true
    deb-systemd-helper enable aptupdatechecker_firmware_update_checker.timer || true
    deb-systemd-helper enable aptupdatechecker_autoclean.timer || true
    deb-systemd-helper enable aptupdatechecker_autoremove.timer || true

    # Start the main service and timers
    deb-systemd-invoke start aptupdatechecker.service || true
    deb-systemd-invoke start aptupdatechecker_apt_update_checker.timer || true
    deb-systemd-invoke start aptupdatechecker_firmware_update_checker.timer || true
    deb-systemd-invoke start aptupdatechecker_autoclean.timer || true
    deb-systemd-invoke start aptupdatechecker_autoremove.timer || true
fi

exit 0
