# systemd/user/aptupdatechecker_apt_update_checker.path
[Unit]
Description=Watch for APT cache update completion
Documentation=https://github.com/tyler/aptupdatechecker

[Path]
PathExists=%t/apt-updates-available

[Install]
WantedBy=default.target
